# debian/postinst
#!/bin/bash

echo "running oems postinstall"
sudo systemctl enable openems.service
sudo systemctl start openems.service
